# pi5/04_implementation/monitoring/systemd/pi-monitor.service
[Unit]
Description=Pi monitoring checks

[Service]
Type=oneshot
ExecStart=/srv/monitoring/bin/monitor-checks.sh
