# pi5/04_implementation/fail2ban/reports/fail2ban-monthly-report.cron
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
0 8 1 * * root /usr/local/bin/fail2ban-monthly-report.sh
