# pi5/04_implementation/x120x_ups/systemd/x120x-ups-shutdown.service
[Unit]
Description=Geekworm X12xx UPS low-battery shutdown handler
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/x120x-ups-shutdown.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
