# pi5/04_implementation/x120x_ups/systemd/x120x-ups-report.timer
[Unit]
Description=Run monthly UPS battery report

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
