# pi5/04_implementation/x120x_ups/systemd/x120x-ups-report.service
[Unit]
Description=Monthly UPS battery report email

[Service]
Type=oneshot
ExecStart=/usr/local/bin/x120x-ups-report.py
