# pi5/04_implementation/nextcloud/clamav-mirror/systemd/clamav-mirror-update.timer
[Unit]
Description=Run cvdupdate hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
