# pi5/04_implementation/nextcloud/clamav-mirror/systemd/clamav-mirror-update.service
[Unit]
Description=Update ClamAV local mirror (cvdupdate)

[Service]
Type=oneshot
User=xindy
WorkingDirectory=/home/xindy
ExecStart=/usr/bin/cvdupdate update
