# pi5/04_implementation/smart/systemd/smart-short-test.service
[Unit]
Description=Daily SMART short test

[Service]
Type=oneshot
ExecStart=/usr/sbin/smartctl -t short /dev/sda
